# complements/fs/mail/etc/local.d/msmtp-sendmail.start
#!/bin/sh
ln -sf /usr/bin/msmtp /usr/bin/sendmail
ln -sf /usr/bin/msmtp /usr/sbin/sendmail